# packaging/update-base.service
[Unit]
Description=Update netboot base image with security patches
After=network.target nfs-server.service

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash -c 'for d in /exports/netboot/*/; do [ -d "$$d/usr" ] && systemd-nspawn -D "$$d" apt-get update && systemd-nspawn -D "$$d" apt-get upgrade -y && systemd-nspawn -D "$$d" apt-get autoremove -y; done'
StandardOutput=journal
StandardError=journal
